# prc.sh
#!/bin/bash
#chkconfig: 2345 80 90
#description:test
nohup /bin/bash /root/lk.sh >> /root/log &
exit
